# ansible/ray_env/ray-worker.service
[Unit]
Description=Ray Worker Node
After=network.target

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/home/ubuntu/ray_worker
ExecStart=/home/ubuntu/ray_worker/start_ray.sh
Restart=on-failure
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target